# osint_app/osint-ingestion/config/argos-ingestion.timer
[Unit]
Description=Run ARGOS Intelligent Ingestion every 30 minutes
Requires=argos-ingestion.service

[Timer]
# Run every 30 minutes
OnCalendar=*:0/30
# Run immediately on startup
OnBootSec=5min
# If the system was down when it should have run, run it when it starts up
Persistent=true

[Install]
WantedBy=timers.target